# usr/lib/systemd/system/android-mount.service.d/10-yggdrasilx-vendor.conf
[Service]
ExecStartPre=mount /dev/disk/by-partlabel/vendor /vendor -t ext4 -o ro
